# files/rc.local
#!/bin/bash
#
# Start Minecraft
#
screen -dmS minecraft bash -c "cd /usr/local/minecraft && java -jar server.jar"
