# files/rc.local
#!/bin/bash
#
# Start Minecraft
#
screen -dmS minecraft bash -c "sleep 10 ; cd /usr/local/minecraft ; java -jar server.jar"
